# .config/systemd/user/udev.service
[Unit]
Description=start udev

[Service]
Type=oneshot
ExecStart=sudo udevadm trigger
