# .config/autostart/slack.desktop
[Desktop Entry]
Type=Application
Comment=slack
Name=slack
Exec=flatpak run com.slack.Slack
Terminal=false
Icon=wine
